# tools/PwnageTool.app/Contents/Resources/CustomPackages/CydiaInstaller.bundle/files/private/var/lib/dpkg/info/cydia.list
/.
/usr
/usr/libexec
/usr/libexec/cydia
/usr/libexec/cydia/space.sh
/usr/libexec/cydia/firmware.sh
/usr/libexec/cydia/finish.sh
/usr/libexec/cydia/startup
/usr/libexec/cydia/free.sh
/usr/libexec/cydia/move.sh
/usr/libexec/cydia/asuser
/usr/libexec/cydia/du
/Applications
/Applications/Cydia.app
/Applications/Cydia.app/commercial.png
/Applications/Cydia.app/Modes
/Applications/Cydia.app/Modes/REMOVE.png
/Applications/Cydia.app/Modes/INSTALL.png
/Applications/Cydia.app/Modes/NEW_INSTALL.png
/Applications/Cydia.app/Modes/PURGE.png
/Applications/Cydia.app/Modes/UPGRADE.png
/Applications/Cydia.app/Modes/REINSTALL.png
/Applications/Cydia.app/Modes/DOWNGRADE.png
/Applications/Cydia.app/installed.png
/Applications/Cydia.app/confirm.html
/Applications/Cydia.app/Test
/Applications/Cydia.app/configure.png
/Applications/Cydia.app/package.html
/Applications/Cydia.app/Info.plist
/Applications/Cydia.app/expanded.png
/Applications/Cydia.app/folder.png
/Applications/Cydia.app/reload.png
/Applications/Cydia.app/signature.html
/Applications/Cydia.app/German.lproj
/Applications/Cydia.app/German.lproj/Sections.strings
/Applications/Cydia.app/German.lproj/Localizable.strings
/Applications/Cydia.app/Default.png
/Applications/Cydia.app/web.png
/Applications/Cydia.app/home@2x.png
/Applications/Cydia.app/folder@2x.png
/Applications/Cydia.app/search.png
/Applications/Cydia.app/install.png
/Applications/Cydia.app/Icon-Small.png
/Applications/Cydia.app/Cydia_
/Applications/Cydia.app/unknown.html
/Applications/Cydia.app/console.png
/Applications/Cydia.app/filesystem.png
/Applications/Cydia.app/installing.png
/Applications/Cydia.app/unknown.png
/Applications/Cydia.app/Sections.plist
/Applications/Cydia.app/trusted.png
/Applications/Cydia.app/sources.png
/Applications/Cydia.app/Sections
/Applications/Cydia.app/Sections/Multimedia.png
/Applications/Cydia.app/Sections/Security.png
/Applications/Cydia.app/Sections/X Window.png
/Applications/Cydia.app/Sections/Packaging.png
/Applications/Cydia.app/Sections/Carrier Bundles.png
/Applications/Cydia.app/Sections/Scripting.png
/Applications/Cydia.app/Sections/Health and Fitness.png
/Applications/Cydia.app/Sections/Networking.png
/Applications/Cydia.app/Sections/Ringtones.png
/Applications/Cydia.app/Sections/Planet-iPhones Mods.png
/Applications/Cydia.app/Sections/Productivity.png
/Applications/Cydia.app/Sections/Tweaks.png
/Applications/Cydia.app/Sections/Dictionaries.png
/Applications/Cydia.app/Sections/Data Storage.png
/Applications/Cydia.app/Sections/Social.png
/Applications/Cydia.app/Sections/SBSettings Addons.png
/Applications/Cydia.app/Sections/WebClips.png
/Applications/Cydia.app/Sections/Communication.png
/Applications/Cydia.app/Sections/Blanks.png
/Applications/Cydia.app/Sections/Education.png
/Applications/Cydia.app/Sections/Utilities.png
/Applications/Cydia.app/Sections/Games.png
/Applications/Cydia.app/Sections/Entertainment.png
/Applications/Cydia.app/Sections/Messaging.png
/Applications/Cydia.app/Sections/Repositories.png
/Applications/Cydia.app/Sections/Navigation.png
/Applications/Cydia.app/Sections/Widgets.png
/Applications/Cydia.app/Sections/Toys.png
/Applications/Cydia.app/Sections/System.png
/Applications/Cydia.app/Sections/Development.png
/Applications/Cydia.app/Sections/Archiving.png
/Applications/Cydia.app/Sections/Wallpaper.png
/Applications/Cydia.app/Sections/App Addons.png
/Applications/Cydia.app/Sections/Keyboards.png
/Applications/Cydia.app/Sections/Text Editors.png
/Applications/Cydia.app/Sections/Imaging.png
/Applications/Cydia.app/Sections/Localization.png
/Applications/Cydia.app/Sections/eBooks.png
/Applications/Cydia.app/Sections/Java.png
/Applications/Cydia.app/Sections/Themes.png
/Applications/Cydia.app/Sections/Administration.png
/Applications/Cydia.app/Sections/Terminal Support.png
/Applications/Cydia.app/Purposes
/Applications/Cydia.app/Purposes/commercial.png
/Applications/Cydia.app/Purposes/extension.png
/Applications/Cydia.app/Purposes/library.png
/Applications/Cydia.app/Purposes/console.png
/Applications/Cydia.app/Purposes/daemon.png
/Applications/Cydia.app/Purposes/uikit.png
/Applications/Cydia.app/Purposes/x.png
/Applications/Cydia.app/storage.html
/Applications/Cydia.app/Cydia
/Applications/Cydia.app/green-dn.png
/Applications/Cydia.app/storage.png
/Applications/Cydia.app/Sources
/Applications/Cydia.app/Sources/cy.sosiphone.com.png
/Applications/Cydia.app/Sources/www.zodttd.com.png
/Applications/Cydia.app/Sources/iphone.hackndev.org.png
/Applications/Cydia.app/Sources/iphonehe.com.png
/Applications/Cydia.app/Sources/urbanfanatics.com.png
/Applications/Cydia.app/version.png
/Applications/Cydia.app/packages.png
/Applications/Cydia.app/green-up.png
/Applications/Cydia.app/icon@2x.png
/Applications/Cydia.app/manage.html
/Applications/Cydia.app/search@2x.png
/Applications/Cydia.app/French.lproj
/Applications/Cydia.app/French.lproj/Sections.strings
/Applications/Cydia.app/French.lproj/Localizable.strings
/Applications/Cydia.app/icon-72.png
/Applications/Cydia.app/changes@2x.png
/Applications/Cydia.app/menes
/Applications/Cydia.app/menes/style.css
/Applications/Cydia.app/menes/chevron.png
/Applications/Cydia.app/menes/chevron@2x.png
/Applications/Cydia.app/menes/indicator-c7ced5-4d4d70.gif
/Applications/Cydia.app/menes/menes.js
/Applications/Cydia.app/home-Selected.png
/Applications/Cydia.app/package.js
/Applications/Cydia.app/chevron@2x.png
/Applications/Cydia.app/home.png
/Applications/Cydia.app/install@2x.png
/Applications/Cydia.app/manage@2x.png
/Applications/Cydia.app/storage.js
/Applications/Cydia.app/settings.png
/Applications/Cydia.app/email.png
/Applications/Cydia.app/Default-Portrait.png
/Applications/Cydia.app/changes.png
/Applications/Cydia.app/confirm.js
/Applications/Cydia.app/installed@2x.png
/Applications/Cydia.app/star.png
/Applications/Cydia.app/English.lproj
/Applications/Cydia.app/English.lproj/Sections_.strings
/Applications/Cydia.app/English.lproj/Localizable.strings
/Applications/Cydia.app/error.html
/Applications/Cydia.app/compose.png
/Applications/Cydia.app/advanced.png
/Applications/Cydia.app/loading.html
/Applications/Cydia.app/home-Selected@2x.png
/Applications/Cydia.app/manage.png
/Applications/Cydia.app/localize.js
/Applications/Cydia.app/removing.png
/Applications/Cydia.app/id.png
/Applications/Cydia.app/warning.png
/Applications/Cydia.app/icon.png
/var
/var/lib
/var/lib/cydia
/System
/System/Library
/System/Library/LaunchDaemons
/System/Library/LaunchDaemons/com.saurik.Cydia.Startup.plist
/System/Library/PreferenceBundles
/System/Library/PreferenceBundles/CydiaSettings.bundle
/System/Library/PreferenceBundles/CydiaSettings.bundle/icon.png
