# tools/PwnageTool.app/Contents/Resources/CustomPackages/CydiaInstallerATV.bundle/files/private/var/lib/dpkg/info/openssh.list
/.
/usr
/usr/bin
/usr/bin/ssh-agent
/usr/bin/slogin
/usr/bin/ssh-add
/usr/bin/ssh-keygen
/usr/bin/ssh-keyscan
/usr/bin/scp
/usr/bin/ssh
/usr/bin/sftp
/usr/share
/usr/share/Ssh.bin
/usr/sbin
/usr/sbin/sshd
/usr/libexec
/usr/libexec/sshd-keygen-wrapper
/usr/libexec/ssh-keysign
/usr/libexec/ssh-rand-helper
/usr/libexec/sftp-server
/Library
/Library/LaunchDaemons
/Library/LaunchDaemons/com.openssh.sshd.plist
/var
/var/empty
/etc
/etc/ssh
/etc/ssh/ssh_config
/etc/ssh/sshd_config
/etc/ssh/moduli
